# com.myron.meetscribe.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.myron.meetscribe</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/myron/projects/meetscribe/scripts/watch-handler.sh</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>/Users/myron/Videos/OBS</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>LowPriorityIO</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/myron/projects/meetscribe/.logs/pipeline.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/myron/projects/meetscribe/.logs/pipeline.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/opt/ffmpeg/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/myron</string>
    </dict>
</dict>
</plist>
